# services/teela.service
[Unit]
Description=Teela Humanoid Robot Brain
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=roni
WorkingDirectory=/home/roni/teela_brain_5
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="KIMI_API_KEY=%KIMI_API_KEY%"
ExecStartPre=/usr/bin/python3 -c "import cv2; print('cv2 ok')"
ExecStart=/usr/bin/python3 -m scripts.conversation_loop --config config.yaml
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
